# trunk/bin/com.apple.macbe.createUser_seoscan.plist
<?xml version="1.0" encoding="UTF-8"?> 
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" 
"http://www.apple.com/DTDs/PropertyList-1.0.dtd"> 
<plist version="1.0"> 
<dict> 
        <key>Label</key> 
        <string>com.apple.macbe.createUser_seoscan.plist</string> 
        <key>ProgramArguments</key> 
        <array> 
                <string>/usr/bin/nice</string>
                <string>/usr/sbin/createUser_seoscan</string> 
                <string>deleteafter</string> 
        </array> 
        <key>RunAtLoad</key> 
        <true/>
        <key>ServiceDescription</key>
        <string>create user seoscan</string>
        <key>StandardErrorPath</key>
        <string>/var/log/com.apple.macbe.createUser_seoscan.log</string>
        <key>StandardOutPath</key>
        <string>/var/log/com.apple.macbe.createUser_seoscan.log</string>
</dict> 
</plist>
